# systemd/refresh_gh.timer
[Unit]
Description=Mise à jour de la liste des GHs

[Timer]
OnCalendar=*:0/5
Persistent=false
AccuracySec=1s

[Install]
WantedBy=timers.target